# telegram-bot.service
[Unit]
Description=Telegram Corporate Bot Service
After=network.target

[Service]
Type=simple
User=mixa-root
WorkingDirectory=/Users/mixa-root/Desktop/бот по работе
Environment=PATH=/Users/mixa-root/Desktop/бот по работе/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/Users/mixa-root/Desktop/бот по работе/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 